# systemd/shotdeck-update-response.path
[Unit]
Description=Watch for Shotdeck update prompt responses

[Path]
PathExists=/opt/shotdeck/update-response.json
PathChanged=/opt/shotdeck/update-response.json
Unit=shotdeck-update-response.service

[Install]
WantedBy=multi-user.target
