# systemd/shotdeck-updater.service
[Unit]
Description=Shotdeck updater
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/shotdeck/updater.json

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/shotdeck
EnvironmentFile=-/etc/default/shotdeck
EnvironmentFile=-/etc/default/shotdeck-updater
ExecCondition=/usr/bin/python3 /opt/shotdeck/public-updater-installer/scripts/check_update_connectivity.py /etc/shotdeck/updater.json
ExecStart=/bin/bash /opt/shotdeck/public-updater-installer/scripts/run_configured_updater.sh /etc/shotdeck/updater.json scheduled
StandardOutput=journal
StandardError=journal
